# upbit-trader/deploy/rebalance.timer
[Unit]
Description=Run portfolio rebalance monthly

[Timer]
# 매월 1일 00:10(서버시간) 실행 + 부팅 직후 한 번(놓친 주기 만회)
OnCalendar=*-*-01 00:10:00
OnBootSec=3min
Persistent=true
Unit=rebalance.service

[Install]
WantedBy=timers.target
